# supervisord.conf
[supervisord]
logfile=/data/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
childlogdir=/data/logs
nodaemon=true
user=root

[inet_http_server]
port=*:9001

[program:api]
directory=/app
command=uvicorn src.api.app:app --workers 1 --host 0.0.0.0 --port 8000 --log-config /app/src/api/uvicorn_config.json
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
killasgroup=true
stopasgroup=true
user=app
environment=PYTHONPATH=/app


[program:aps]
directory=/app
command=python -m src.start.commands aps
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
killasgroup=true
user=app
stopasgroup=true
environment=PYTHONPATH=/app
